# files/etc/systemd/system/prometheus.service
[Unit]
Description=Prometheus Server
Documentation=https://prometheus.io/docs/introduction/overview/
After=network-online.target

[Service]
User=prometheus
Group=prometheus
Restart=on-failure
ProtectSystem=strict
ProtectHome=true
TimeoutSec=10s
ReadWritePaths=/var/lib/prometheus
ExecStart=/opt/prometheus/prometheus --storage.tsdb.path=/var/lib/prometheus --storage.tsdb.retention.time=1d --config.file=/etc/prometheus/config.yml

[Install]
WantedBy=multi-user.target
